# demo/net/epoll-tcp-server/epoll-tcp-server.service
[Unit]
Description=TCP epoll server

[Service]
ExecStart=/opt/epoll-tcp-server/epoll-tcp-server 172.16.6.17

[Install]
WantedBy=multi-user.target
