# net/tcp-server/tcp-server.service
[Unit]
Description=TCP server

[Service]
ExecStart=/opt/tcp-server/tcp-server 172.16.6.17

[Install]
WantedBy=multi-user.target
